# etc/auksdrenewer.service.in
[Unit]
Description=Auks Credentials Renewer Daemon
After=network.target
Wants=aukspriv.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/auksdrenewer
ExecStart=@sbindir@/auksdrenewer -F $AUKSDRENEWER_OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
