# etc/aukspriv.service.in
[Unit]
Description=Auks ccache from keytab scripted daemon
After=network.target

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/aukspriv
ExecStart=@sbindir@/aukspriv $AUKSPRIV_OPTIONS
GuessMainPID=true

[Install]
WantedBy=multi-user.target
